# roastery/etc/systemd/system/var-mnt-data.automount
[Unit]
Description=Automount data directory

[Automount]
TimeoutIdleSec=20min
Where=/var/mnt/data

[Install]
WantedBy=multi-user.target
