# barista/rootfs/etc/systemd/system/populate-secret@.service
[Unit]
Description=Fetch and Register Podman Secret
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/populate-secret %i
