# barista/rootfs/etc/systemd/system/populate-secret@.service
[Unit]
Description=Fetch and Register Podman Secret

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/populate-secret %i
